# SmartCarNetCore/SmartCar.Web/smartcarwebserver.service
[Unit]
Description=Run NET Server


[Service]
ExecStart=/home/ubuntu/SmartCar/SmartCarNetCore/SmartCar.Web/run.sh

[Install]
WantedBy=default.target